# roles/afraid_ddns/templates/ddns_update.service.j2
[Unit]
Description=Job to update DDNS entry

[Service]
Type=simple
ExecStart=/home/{{ ansible_user }}/ddns_update/ddns_update.sh

[Install]
WantedBy=default.target
